# systemd/kidbot-updater.timer
[Unit]
Description=Manual-only KidBot updater timer (disabled by install.sh)

[Timer]
OnBootSec=2min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
